# deploy/shiolplus.service
[Unit]
Description=SHIOL+ FastAPI Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/SHIOL-PLUS
Environment="PATH=/root/.venv_shiolplus/bin"
ExecStart=/root/.venv_shiolplus/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

# Graceful shutdown configuration
# Give the process 180 seconds (3 minutes) to terminate gracefully before SIGKILL
# This allows long-running polling operations to clean up properly
TimeoutStopSec=180
KillMode=mixed
KillSignal=SIGTERM

# Resource limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
